# board/uconsole/overlay/etc/network/if-pre-up.d/wpa_supplicant
#!/bin/sh

if [ "${IFACE}" == "wlan0" ]; then
	wpa_supplicant -B -i ${IFACE} -c /etc/wpa_supplicant.conf
fi
